# infra/orchestrator.service
[Unit]
Description=RGrid Orchestrator Daemon
After=network.target postgresql.service

[Service]
Type=simple
User=rgrid
WorkingDirectory=/opt/rgrid
EnvironmentFile=/opt/rgrid/.env
ExecStart=/opt/rgrid/venv/bin/python -m orchestrator.daemon
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rgrid-orchestrator

[Install]
WantedBy=multi-user.target
